# deploy/biorxiv-sync.service.in
[Unit]
Description=bioRxiv database sync

[Service]
Type=oneshot
User=@RUN_USER@
Group=@RUN_USER@
WorkingDirectory=@PROJECT_ROOT@
EnvironmentFile=-@PROJECT_ROOT@/deploy/biorxiv-mcp.env
ExecStart=@PROJECT_ROOT@/.venv/bin/biorxiv-mcp-sync
# Bulk sync can take hours
TimeoutStartSec=86400
Restart=on-failure
RestartSec=30
