# deploy/biorxiv-mcp.service
[Unit]
Description=bioRxiv MCP server
After=network.target

[Service]
WorkingDirectory=/home/hmblair/dev/biorxiv-mcp
ExecStart=/home/hmblair/dev/biorxiv-mcp/.venv/bin/biorxiv-mcp
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
