# deploy/biorxiv-sync.service
[Unit]
Description=bioRxiv database sync

[Service]
Type=oneshot
WorkingDirectory=/home/hmblair/dev/biorxiv-mcp
ExecStart=/home/hmblair/dev/biorxiv-mcp/.venv/bin/biorxiv-mcp-sync
# Bulk sync can take hours
TimeoutStartSec=86400
Restart=on-failure
RestartSec=30
