# deploy/biorxiv-sync.service.in
[Unit]
Description=bioRxiv database sync

[Service]
Type=oneshot
WorkingDirectory=@PROJECT_ROOT@
ExecStart=@PROJECT_ROOT@/.venv/bin/biorxiv-mcp-sync
# Bulk sync can take hours
TimeoutStartSec=86400
Restart=on-failure
RestartSec=30
